# systemd/salessync-startup.service
[Unit]
Description=SalesSync Application Startup Service
Documentation=https://github.com/Reshigan/SalesSyncAI
After=network.target postgresql.service nginx.service
Wants=postgresql.service nginx.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/SalesSyncAI
RemainAfterExit=yes

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=DATABASE_URL=postgresql://salessync:salessync123@localhost:5432/salessync
Environment=JWT_SECRET=your-super-secret-jwt-key-change-in-production
Environment=BCRYPT_ROUNDS=12

# Startup script
ExecStart=/home/ubuntu/SalesSyncAI/scripts/startup.sh

# Restart policy
Restart=on-failure
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=salessync-startup

[Install]
WantedBy=multi-user.target